# vllm-server.service
[Unit]
Description=vLLM Server on Port 5002
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/SLYDInference
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /root/SLYDInference/launch_vllm_server.py --config /root/SLYDInference/vllm_config_port5002.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target